# crontab
@reboot	/home/pi/MagicMirror/simplehttp.sh
*/2 * * * * /home/pi/MagicMirror/getics.sh
@reboot	/home/pi/MagicMirror/run.sh
